# systemd/user/ping_log.timer
[Unit]
Description=Run "ping_log.service" every 30 seconds

[Timer]
OnBootSec=30
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
