# systemd/shurjo/keep_audio_alive.timer
[Unit]
Description=Run "keep_alive_audio.service] every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
